# systemd/lxmd.service
[Unit]
Description=LXMF Propagation Node Daemon
After=network.target rnsd.service
Wants=network.target
Requires=rnsd.service

[Service]
Type=simple
User=rns
Group=rns
ExecStart=/home/rns/.local/bin/lxmd --propagation-node --rnsconfig /etc/reticulum
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
