# systemd/rrcd.service
[Unit]
Description=Reticulum Relay Chat Daemon
After=network.target rnsd.service
Requires=rnsd.service

[Service]
Type=simple
User=rns
Group=rns
Environment=PATH=/home/rns/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/rns
ExecStart=/home/rns/.local/bin/rrcd
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
